# cron/job.cron
* * * * * echo "Cron running at $(date -u)" >> /data/cron.log
